# packages/asteria-openrgb/asteria_openrgb-0.1.0.tar.gz/asteria_openrgb-0.1.0/asteria.service
[Unit]
Description=Framework to automatically control LEDs via OpenRGB

[Service]
WorkingDirectory=/home/josh/projects/asteria
ExecStart=/home/josh/projects/asteria/test-venv/bin/asteria config.toml

[Install]
WantedBy=multi-user.target
